# systemd/teamtp-support-bot.service
[Unit]
Description=TeamTP Support Bot
After=teamspeak6.service
Wants=teamspeak6.service

[Service]
Type=simple
User=teamtp
Group=teamtp
WorkingDirectory=/opt/teamtp/bots/support-bot
ExecStart=/usr/bin/node /opt/teamtp/bots/support-bot/index.js
EnvironmentFile=/opt/teamtp/.env
Restart=on-failure
RestartSec=10
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
