# systemd/teamspeak6.service
[Unit]
Description=TeamSpeak 6 Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=tsserver
Group=tsserver
WorkingDirectory=/opt/teamtp/server/teamspeak6
ExecStart=/opt/teamtp/server/teamspeak6/tsserver --config-file /opt/teamtp/server/teamspeak6/tsserver.yaml
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
LimitNPROC=4096
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
