# systemd/teamtp-panel.service
[Unit]
Description=TeamTP Web Panel
After=teamspeak3.service
Wants=teamspeak3.service

[Service]
Type=simple
User=teamtp
Group=teamtp
WorkingDirectory=/opt/teamtp/panel
ExecStart=/usr/bin/node /opt/teamtp/panel/server.js
EnvironmentFile=/opt/teamtp/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
